# djsupervisor/contrib/django_ztask/supervisord.conf
;
;  Basic supervisord.conf for django-ztask.
;
;  You get a program named "ztaskd".
;


[program:ztaskd]

;  Show extra logging output if running in debug mode.
{% if settings.DEBUG %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py ztaskd --noreload --loglevel=INFO
{% else %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py ztaskd --noreload
{% endif %}

directory={{ PROJECT_DIR }}

