# djsupervisor/contrib/django_wsgiserver/supervisord.conf
;
;  Basic supervisord.conf for django-wsgiserver.
;
;  If django_cpserver is installed, we override the default runserver
;  command to use runcpserver instead.  That's almost always what you want.
;

[program:runserver]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runwsgiserver

