# djsupervisor/contrib/django_wsgiserver/supervisord.conf

;  If django_wsgiserver is installed, override the default
;  runserver command to use runwsgiserver instead.

[program:runserver]
command={{ PROJECT_DIR }}/manage.py runwsgiserver

