# djsupervisor/contrib/djcelery/supervisord.conf
;
;  Basic supervisord.conf for django-celery.
;
;  This is based on the file distributed with django-celery, but
;  modified for our purposes here.
;
;  You get two programs.  "celeryd" is the main task-handling daemon,
;  and "celerybeat" is the periodic task scheduler.
;


[program:celeryd]

;  Show extra logging output if running in debug mode.
{% if settings.DEBUG %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd --loglevel=INFO
{% else %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd
{% endif %}

directory={{ PROJECT_DIR }}
numprocs=1
startsecs=5

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998


[program:celerybeat]

{% if settings.DEBUG %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celerybeat --loglevel=INFO
{% else %}
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celerybeat
{% endif %}

directory={{ PROJECT_DIR }}
numprocs=1
startsecs=5

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999

