# djsupervisor/contrib/djcelery/supervisord.conf

[program:celeryd]
command={{ PROJECT_DIR }}/manage.py celeryd

