# djsupervisor/contrib/django_cpserver/supervisord.conf

;  If django_cpserver is installed, override the default
;  runserver command to use runcpserver instead.
[program:runserver]
command={{ PROJECT_DIR }}/manage.py runcpserver

